# launchd/local.backup-tools.full.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Weekly full backup + SHA-256 hash verify (Sunday). Installed by install.sh,
     which substitutes @RUN_SH@ / @LOG@ / @PATH@. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.backup-tools.full</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@RUN_SH@</string>
  </array>

  <!-- Sundays at 03:00 (Weekday 0 = Sunday). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>@PATH@</string>
  </dict>

  <key>StandardOutPath</key>
  <string>@LOG@</string>
  <key>StandardErrorPath</key>
  <string>@LOG@</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
